# systemd/nws-alerts.service
[Unit]
Description=NWS alerts thermal receipt printer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jamesposs
WorkingDirectory=/opt/nws-alerts
EnvironmentFile=/opt/nws-alerts/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/nws-alerts/.venv/bin/python -m app.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
